# docker/cron/crontab
* * * * * www-data /usr/local/bin/php /var/www/symfony/bin/console app:some-command >> /var/log/cron.log 2>&1
